# systemd/wildlife-analysis.timer
[Unit]
Description=Wildlife Monitor — Daily Analysis Timer
Requires=wildlife-analysis.service

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
